# example/trojan-go.service
[Unit]
Description=Trojan-Go Service (Modern Super Suite)
After=network.target nss-lookup.target
Documentation=https://github.com/p4gefau1t/trojan-go

[Service]
# 环境建议：将配置文件放在 /etc/trojan-go/ 目录下
# 将二进制文件放在 /usr/local/bin/ 目录下
User=root
# 赋予监听 443 等特权端口的权限（即使非 Root 用户运行也建议设置）
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# 启动命令：指向 YAML 配置文件
ExecStart=/usr/local/bin/trojan-go -config /etc/trojan-go/server.yaml

Restart=on-failure
RestartSec=10s
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target
